# config/apsta.service
[Unit]
Description=AP STA mode service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=setup_networking

[Install]
WantedBy=multi-user.target