# config/apsta.service
[Unit]
Description=AP STA mode service

[Service]
Type=oneshot
ExecStart=startDualNetworkMode

[Install]
WantedBy=sysinit.target